# docs/reference-architecture/cloud-based-document-management-system/guides/implementing-cloud-based-document-management-system/cbdms-terraform/documents/cbdms_unison.service
[Unit]
Description=Sync NGINX Certificates
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/unison cbdms-app-unison

[Install]
WantedBy=multi-user.target

