# docs/reference-architecture/cloud-based-document-management-system/guides/implementing-cloud-based-document-management-system/cbdms-terraform/documents/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=multi-user.target

